# scripts/systemd/services/clear-trash.service
[Unit]
Description=Clears trash bin

[Service]
ExecStart=rm -rf /drives/data/.trash/*

[Install]
WantedBy=multi-user.target
