# scripts/systemd/timers/ntpdate.timer
[Unit]
Description=Starts ntp after boot

[Timer]
OnBootSec=60s
Unit=ntpdate.service

[Install]
WantedBy=timers.target
